# CM4 Waveshare Dual DNS LCD/systemd/cm4-safe-nightly-update.service
[Unit]
Description=CM4 safe Bookworm-only nightly apt update and reboot
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cm4-safe-nightly-update
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
